# examples/pvemcp.service
[Unit]
Description=pvemcp MCP server
After=network.target

[Service]
Type=simple
User=root
Environment=PVEMCP_AUDIT_LOG=/var/log/pvemcp-audit.log
Environment=PVEMCP_ALLOW_BREAK_GLASS=1
ExecStart=pvemcp-server
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
